# supervisor.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:civic_tracker]
command=gunicorn -c gunicorn.conf.py app:app
directory=/app
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/civic_tracker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:scheduler]
command=python -c "from app import setup_action_updates; setup_action_updates()"
directory=/app
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/scheduler.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10

[program:backup_scheduler]
command=python -c "from apscheduler.schedulers.blocking import BlockingScheduler; from scripts.backup import BackupManager; sched = BlockingScheduler(); sched.add_job(BackupManager().run_backup, 'cron', day_of_week='sun', hour=0, minute=0); sched.start()"
directory=/app
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/backup.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
